# pdf-sign-sdk.service
[Unit]
Description=PDF Sign SDK Web Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/pdf-sign-sdk
ExecStart=/usr/bin/node /var/www/pdf-sign-sdk/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pdf-sign-sdk

# Environment
Environment=NODE_ENV=production
Environment=PORT=8000
Environment=HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
